# deploy/aiclass-frontend.service
# ============================================
# AI不难学 - 前端 Nuxt 3 SSR 服务
# 放到 /etc/systemd/system/aiclass-frontend.service
# ============================================

[Unit]
Description=AI不难学 Frontend (Nuxt 3)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aiclass/frontend
Environment="NODE_ENV=production"
Environment="NITRO_PORT=3000"
Environment="NITRO_HOST=127.0.0.1"
ExecStart=/usr/bin/node .output/server/index.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
